# thinkpad-x1e-config/etc/systemd/system/windows.automount
[Unit]
Description=Automount Windows C:\ Drive

[Automount]
Where=/windows

[Install]
WantedBy=multi-user.target
