# mach-thinkpad-x1e/etc/systemd/system/boot.automount
[Unit]
Description=Automount UEFI partition
ConditionVirtualization=False

[Automount]
Where=/boot

[Install]
WantedBy=multi-user.target
