# mach-thinkpad-x1e/etc/systemd/system/boot.automount
[Unit]
Description=Automount UEFI partition

[Automount]
Where=/boot

[Install]
WantedBy=multi-user.target
